# .config/systemd/user/Mailspring.service
[Unit]
Description=Starts Mailspring

[Service]
ExecStart=mailspring
