# .config/systemd/user/numlockx.service
[Unit]
Description=Turns numlock on

[Service]
ExecStart=numlockx on
