# .config/systemd/user/polkit-gnome-authentication-agent-1.service
[Unit]
Description=Starts GNOME polkit auth agent

[Service]
ExecStart=/usr/libexec/polkit-gnome-authentication-agent-1
