# .config/systemd/user/Steam.service
[Unit]
Description=Starts Steam

[Service]
ExecStart=steam -silent
